# target/linux/mediatek/mt7622/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

boot() {
	case $(board_name) in
	dlink,eagle-pro-ai-m32-a1|\
	dlink,eagle-pro-ai-r32-a1)
		;;
	linksys,e8450)
		mtd erase senv || true
		;;
	esac
}
